# systemd/calendar.service
[Unit]
Description=Unusual Whales Calendar Feed
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=ezb0t
Group=ezb0t
WorkingDirectory=/home/ezb0t/ezbot/data_feeds/calendar_feed
ExecStart=/usr/bin/python3 /home/ezb0t/ezbot/data_feeds/calendar_feed/uw_calendar.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Environment
Environment=PYTHONPATH=/home/ezb0t/ezbot
Environment=PYTHONUNBUFFERED=1

# Security
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/home/ezb0t/ezbot

[Install]
WantedBy=multi-user.target